# deploy/linux/message-service.service
[Unit]
Description=Message Service — ETL pipeline report aggregator
Documentation=file:///opt/message-service/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=message-service
Group=message-service
WorkingDirectory=/opt/message-service

# Environment variables.
#
# - The inline `Environment=` below pins the path argparse falls
#   back to when `--config` is absent. The value mirrors the
#   ExecStart= argument so both code paths agree.
# - The optional `EnvironmentFile=-` (the leading hyphen makes the
#   file optional; systemd will not fail-to-start if the path is
#   absent) lets operators drop credentials and per-host overrides
#   into a sibling file without editing this unit. Operators
#   typically place SMTP credentials, OAuth tokens, and any
#   secret-bearing values in the env-file rather than in the TOML
#   config, which is committed to the deployment image.
Environment="MESSAGE_SERVICE_CONFIG=/etc/message-service/config.toml"
EnvironmentFile=-/etc/message-service/message-service.env

ExecStart=/opt/message-service/.venv/bin/message-service --config /etc/message-service/config.toml
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillSignal=SIGTERM

# Sandboxing (adjust as needed for deployment)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/message-service /var/log/message-service

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=message-service

[Install]
WantedBy=multi-user.target
